# fs-source.service
[Unit]
Description=FS Source Daemon - OBS Face Detection Switcher (Always-On)
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/venv/bin/python __INSTALL_DIR__/fs_source_daemon.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="DISPLAY=:0"

# Restart conditions
StartLimitBurst=5
StartLimitIntervalSec=60

[Install]
WantedBy=multi-user.target
